# crontab
*/3 * * * * /usr/local/bin/python /usr/app/main.py 
